# Working-ccdc2024/old/Summer_2022/Practice_Scenarios/P1/Developer/crontab
*/2 * * * * root /home/innocentdev/fileTransfer.sh >> /home/innocentdev/cron.log 2>&1
